# ansible/roles/services/postfix/files/courier-imap-ssl.restart.cron
#!/bin/sh
service courier-imap-ssl restart
